# docker/india-secondary-node/india-secondary-node-supervisor.conf
[program:multichain]
command=/bin/bash -c "./tmp/run-blockchain.sh"