# docker/multichain-explorer/multichain-explorer-supervisor.conf
[program:multichain]
command=/bin/bash -c "./tmp/run-blockchain.sh"

[program:multichain-explorer]
directory=multichain-explorer
command=/bin/bash -c "/tmp/./run-multichain-explorer.sh"
